# zoneminder/10-wvdial.conf
[program:lmlvss]
command=/usr/bin/wvdial
autostart=true
autorestart=true
process_name=wvdial
startretries=20
redirect_stderr=true
stdout_logfile=/var/log/wvdial.log

